# stage-2/etc/profile.d/support.sh
/usr/local/bin/support-wrapper &
P="$!"
trap "kill $P ; wait $P" 0
unset P

